# crontab
0 4 * * * /usr/bin/reboot
*/15 * * * * /usr/bin/monitor-proxy > /dev/null
